# ellis/initsql.supervisord.conf
[program:initsql]
command=/usr/share/clearwater/ellis/create_numbers.sh
autostart=true
autorestart=false
startsecs=0
priority=999
