# ralf/ralf.supervisord.conf
[program:ralf]
command=/etc/init.d/ralf start
autostart=true
startretries=10
autorestart=true
