# astaire/rogers.supervisord.conf
[program:rogers]
command=/etc/init.d/rogers run
startretries=100
autostart=true
autorestart=true
