# homer/homer.supervisord.conf
[program:homer]
command=/etc/init.d/homer run
autostart=true
autorestart=true
startretries=5000
