# bono/bono.supervisord.conf
[program:bono]
command=/etc/init.d/bono start
autostart=true
autorestart=true
