# bono/restund.supervisord.conf
[program:restund]
command=/etc/init.d/restund start
autostart=true
autorestart=true
